# contrib/systemd/acme-redirect-renew.service
[Unit]
Description=acme-redirect: renew certs if necessary

[Service]
ExecStart=/usr/bin/acme-redirect renew
