# contrib/systemd/acme-redirect-renew.timer
[Unit]
Description=acme-redirect: renew certs if necessary

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
